# crontab
# This is example crontab record.
# Install it for root

* * * * * /path/to/top_log
